# ansible/roles/tinyauth/templates/tinyauth.service.j2
[Unit]
Description=TinyAuth Authentication Server
After=network.target

[Service]
Type=simple
User=tinyauth
Group=tinyauth
WorkingDirectory=/opt/tinyauth
ExecStart=/opt/tinyauth/bin/tinyauth-start.sh
Restart=always
RestartSec=10
StandardOutput=append:/var/log/tinyauth.log
StandardError=append:/var/log/tinyauth.log

[Install]
WantedBy=multi-user.target
